# docker/supervisor.conf.d/janitoo_audiovideo.conf
[program:jnt_audiovideo]
user=root
autostart=true
autorestart=true
command=jnt_audiovideo -c /opt/janitoo/etc/janitoo_audiovideo.conf start
